# deploy/spielplatzkarte-import.service
[Unit]
Description=Spielplatzkarte OSM data importer
Documentation=https://github.com/mfuhrmann/spielplatzkarte
After=docker.service
Requires=docker.service

[Service]
Type=oneshot

# Set this to the directory where spielplatzkarte is checked out
WorkingDirectory=/opt/spielplatzkarte

# Load OSM_RELATION_ID, PBF_URL, and other env vars from the project .env file
# Update this path to match WorkingDirectory= above
EnvironmentFile=/opt/spielplatzkarte/.env

# Run as the user that owns the deployment directory and is in the docker group
# Uncomment and set to the appropriate user:
#User=spielplatzkarte

ExecStart=docker compose run --rm importer

StandardOutput=journal
StandardError=journal
SyslogIdentifier=spielplatzkarte-import
